# install/ubuntu/upstart.conf
description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

exec /opt/nucleoos/env/bin/uwsgi --ini /opt/nucleoos/uwsgi.ini --die-on-term >> /var/log/nucleoos.log 2>&1
